# docker/php/tasks.cron
# Cron jobs for www-data user
* * * * *   echo "$(date) - cron OK, user - $(whoami)" >> /var/log/cron.log 2>&1
